# server/scripts/com.stockpulse.fundamentals.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.stockpulse.fundamentals</string>

  <!-- node + the one-shot refresh script -->
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/pasupuletisreedhar/Projects/stockpulse/server/scripts/refresh-fundamentals.mjs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/pasupuletisreedhar/Projects/stockpulse/server</string>

  <!-- Run every 6 hours (21600s). If the Mac was asleep, it runs on wake. -->
  <key>StartInterval</key>
  <integer>21600</integer>

  <!-- Also run once shortly after login so data is fresh when you start the day -->
  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/pasupuletisreedhar/Projects/stockpulse/server/scripts/refresh.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/pasupuletisreedhar/Projects/stockpulse/server/scripts/refresh.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
